# clawd-daemon.service
[Unit]
Description=Clawd Usage Dashboard - daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tang0115
WorkingDirectory=/home/tang0115/clawd-dash
# usage.json lives on tmpfs (/dev/shm), not the SD card, to avoid the
# every-15s (every-2s while Spotify plays) write pattern that was traced
# back as the likely cause of repeated SD card corruption. /dev/shm is
# recreated empty on every boot, so the symlink has to be re-pointed here
# rather than committed to git as a static file.
ExecStartPre=/bin/ln -sf /dev/shm/usage.json /home/tang0115/clawd-dash/usage.json
ExecStart=/usr/bin/python3 /home/tang0115/clawd-dash/daemon.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
